# bme688-web.service
[Unit]
Description=BME688 Monitor Web Interface
After=network.target

[Service]
Type=simple
User=luka
WorkingDirectory=/home/luka/bme688-monitor/web
ExecStart=/usr/bin/python3 /home/luka/bme688-monitor/web/server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
